# scripts/templates/teamclaude-relay.service
[Unit]
Description=TeamClaude Relay (injects x-api-key for Claude Code)
After=teamclaude.service
Requires=teamclaude.service

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/teamclaude-relay
# PROXY_API_KEY is generated by server-install-native.sh at first run and
# written to /etc/teamclaude-relay/env (root-only, 0640 group=openclaw).
EnvironmentFile=/etc/teamclaude-relay/env
Environment=TC_HOST=127.0.0.1
Environment=TC_PORT=3456
Environment=PORT=3457
ExecStart=/usr/bin/env node /opt/teamclaude-relay/teamclaude-relay.cjs
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=teamclaude-relay

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
# relay не пишет на диск вообще — выйдем безпривилегий насколько можно
ReadOnlyPaths=/opt/teamclaude-relay
ProtectHome=true

[Install]
WantedBy=multi-user.target
